# docker/marker-supervisord.conf
[program:marker]
user = abacus
directory = /
command = /usr/bin/markerd /data/abacus/marker.conf
startsecs = 2
